# frontend-tier/app.conf
[supervisord]
logfile = /tmp/supervisord.log
user = root
[program:app]
directory=/app
command=gunicorn --worker-class eventlet -w 1 main:app -b 127.0.0.1:8080
autostart=true
autorestart=true
stderr_logfile=/var/log/casapp.err.log
stdout_logfile=/var/log/casapp.out.log
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB